# systemd/user/timers.target.wants/wallpaper.service
[Unit]
Description=Change wallpaper depending on time

[Service]
ExecStart=/home/aj/.config/hypr/scripts/randomPaper.sh



